# TgBotChannelRepost.service
[Service]
WorkingDirectory=/home/????/TgBotChannelRepost
User=????
ExecStart=/usr/bin/python main.py

[Install]
WantedBy=multi-user.target
EOF